# supervisor/armada-version.conf
[program:armada-version]
directory=/opt/armada-version/src
command=python3.5 main.py
